# jffs/crontab/admin
*/5 * * * * /jffs/check_and_swap_wan.sh
